# system/launchagents/com.u29dc.dia-cdp.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- dotfiles-managed: generated dia-cdp -->
<!-- dotfiles-source: system/launchagents/com.u29dc.dia-cdp.plist.template -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.u29dc.dia-cdp</string>

	<key>ProgramArguments</key>
	<array>
		<string>__DIA_APP_BIN__</string>
		<string>--remote-debugging-port=__AGENT_BROWSER_DIA_PORT__</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>ThrottleInterval</key>
	<integer>10</integer>

	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>

	<key>StandardOutPath</key>
	<string>__DIA_STDOUT_LOG__</string>

	<key>StandardErrorPath</key>
	<string>__DIA_STDERR_LOG__</string>
</dict>
</plist>
