# system/launchagents/com.u29dc.dia-cdp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.u29dc.dia-cdp</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Dia.app/Contents/MacOS/Dia</string>
		<string>--remote-debugging-port=9222</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>ThrottleInterval</key>
	<integer>10</integer>

	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>

	<key>StandardOutPath</key>
	<string>/Users/han/Library/Logs/com.u29dc.dia-cdp.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/han/Library/Logs/com.u29dc.dia-cdp.err</string>
</dict>
</plist>
